# overlay/system/usr/share/upstart/sessions/keypad-slide-osk.conf
description "enable/disable OSK on HW keyboard slide event"

start on started dbus
stop on stopped dbus

respawn

exec /usr/bin/keypad-slide-osk
